# calendar-app.service
[Unit]
Description=Touch Calendar Application
After=graphical-session.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=%h/calendar-app
Environment=DISPLAY=:0
ExecStartPre=/bin/sleep 10
ExecStart=%h/calendar-app/start_calendar.sh
Restart=always
RestartSec=3

[Install]
WantedBy=graphical-session.target
